# rdc-webui.service
[Unit]
Description=RDC Jetson Router Web Admin UI
After=rdc-router.service
Wants=rdc-router.service

[Service]
Type=simple
# Ensure we use the Python from the virtual environment
ExecStart=/home/blob/Documents/jetson/router/venv/bin/python /home/blob/Documents/jetson/router/web_app.py
Restart=on-failure
RestartSec=5s
# Change working directory so Flask can find the templates folder correctly
WorkingDirectory=/home/blob/Documents/jetson/router
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
